# evidguard-logwatch/systemd/evidguard-monitor.path
[Unit]
Description=EvidGuard — Vigila el log de rclone manager

[Path]
PathModified=RCLONE_GUI_LOG_PLACEHOLDER

[Install]
WantedBy=multi-user.target
